# install/etc/s6/services/20-php-fpm/run
#!/usr/bin/with-contenv bash

for s in /assets/functions/*; do source $s; done
PROCESS_NAME="php-fpm"

check_container_initialized
check_service_initialized init
liftoff

print_info "Starting php-fpm ${PHP_VERSION}"
silent exec s6-setuidgid nginx php-fpm5 -F;
